# deploy/stenvik-worker.service
[Unit]
Description=Stenvik Agent Studio Worker
After=network.target
Wants=network-online.target

[Service]
Type=simple
User=i48ptgvnis
Group=i48ptgvnis
WorkingDirectory=/home/i48ptgvnis/stenvik-leads
Environment="PATH=/home/i48ptgvnis/stenvik-leads/.venv/bin"
EnvironmentFile=/home/i48ptgvnis/stenvik-leads/.env
ExecStart=/home/i48ptgvnis/stenvik-leads/.venv/bin/python -m worker.main
Restart=always
RestartSec=10
KillSignal=SIGTERM
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=stenvik-worker

# Sandboxing — копия настроек stenvik-leads.service (что точно работает).
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/i48ptgvnis/stenvik-leads/data /home/i48ptgvnis/stenvik-leads/.venv
PrivateTmp=true

[Install]
WantedBy=multi-user.target
